# includes/etc/cron.d/scanner
#
# Default Scanner Scheduled Job
#

*/1 * * * * ansible-playbook nofusscomputing.phpipam_scan_agent.agent -v
