# includes/etc/supervisor/conf.d/update-ca-certificates.conf
[program:trusted-ca-update]
startsecs=0
stopwaitsecs=55
command=/bin/update-ca
autorestart=true
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0